# externals/pd-patch.service
[Unit]
After=mpr121-pd.service

[Service]
ExecStart=/usr/bin/sudo /usr/bin/puredata -nogui /home/pi/patch/main.pd &
Restart=on-failure
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=pd-patch
User=pi
Group=pi
#Environment=NODE_ENV={{ env }}
WorkingDirectory=/home/pi

[Install]
WantedBy=default.target
